# install-fresh-system.d/razer/razer_blackwidow.rules
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="1532", ATTR{idProduct}=="0221", RUN+="/root/init_blackwidow.py"
